# freesia/ansible/roles/adguardhome/files/adguardhome.service
[Unit]
Description=AdGuard Home: Network-level blocker
ConditionFileIsExecutable=/opt/AdGuardHome/AdGuardHome
After=network.target nss-lookup.target

[Service]
StartLimitInterval=5
StartLimitBurst=10
ExecStart=/opt/AdGuardHome/AdGuardHome "-s" "run" "--config" "/opt/AdGuardHome/AdGuardHome.yaml"
Restart=always
RestartSec=10
User=root
EnvironmentFile=-/etc/sysconfig/AdGuardHome
WorkingDirectory=/opt/AdGuardHome

[Install]
WantedBy=multi-user.target
