# freesia/ansible/roles/cardinal/files/cardinal.service
[Unit]
Description=cardinal
[Service]
ExecStart=/usr/local/sbin/cardinal
Type=oneshot
User=root
[Install]
WantedBy=multi-user.target
